# thor/rootfs/usr/lib/systemd/system/systemd-growfs@sysroot.service.d/10-bootc-grow-sysroot.conf
# Grow the root filesystem to fill its partition on first boot.
#
# MountFlags=slave gives the unit a private mount namespace so the remount,rw is
# visible only to systemd-growfs and keeps /sysroot read-only for the host.

[Service]
MountFlags=slave
ExecStartPre=/usr/bin/mount -o remount,rw /sysroot

[Install]
WantedBy=multi-user.target
