# orin/rootfs/usr/lib/systemd/system/ostree-remount.service.d/composefs.conf
# The stock ostree-remount.service requires "ostree" on the kernel
# command line, but bootc composefs uses "composefs=" instead.
[Unit]
ConditionKernelCommandLine=
ConditionKernelCommandLine=composefs
