# thor/rootfs/usr/lib/systemd/system/nvidia-display-load.service
# Bring up the NVIDIA openrm GPU/display stack (Thor / tegra264).
#
# udev autoloads nvidia by PCI modalias when the GPU enumerates, before the
# platform is ready, so that probe fails and the GPU stays unbound. Mirror
# NVIDIA's nv-load-display-modules.sh: unload the early stack and reload it
# during multi-user.target, once the platform is up. (We can't use NVIDIA's
# service directly. Its `depmod -a` writes read-only /usr on bootc, so the
# depmod.d/modprobe.d configs and depmod are baked at build time instead.)
#
[Unit]
Description=Reload NVIDIA openrm GPU/display modules after platform init
After=systemd-modules-load.service
Before=display-manager.service nvpmodel.service

[Service]
Type=oneshot
RemainAfterExit=yes
# Unload the stack in reverse dependency order, ignoring modules that were
# never loaded, then reload it now that the platform is ready.
ExecStartPre=-/usr/sbin/modprobe -r nvidia_uvm
ExecStartPre=-/usr/sbin/modprobe -r nvidia_drm
ExecStartPre=-/usr/sbin/modprobe -r nvidia_modeset
ExecStartPre=-/usr/sbin/modprobe -r nvidia
ExecStart=/usr/sbin/modprobe nvidia_drm

[Install]
WantedBy=multi-user.target
