# orin/rootfs/usr/lib/systemd/system/systemd-growfs@sysroot.service.d/10-bootc-grow-sysroot.conf
# Grow the root filesystem to fill its partition on first boot.
#
# MountFlags=slave gives the unit a private mount namespace so the remount,rw is
# visible only to systemd-growfs and keeps /sysroot read-only for the host. The
# remount and grow must run in one process: systemd sets up the mount namespace
# per Exec* command, so a remount in ExecStartPre would not be seen by ExecStart.
[Service]
MountFlags=slave
ExecStart=
ExecStart=/bin/sh -c '/usr/bin/mount -o remount,rw %f && exec /usr/lib/systemd/systemd-growfs %f'

[Install]
WantedBy=multi-user.target
